# nginx-openstack-heat-demo/nginxplus/sync.conf
description "Sync between backends and NGINX Plus"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid nobody
setgid nogroup

console log

exec python3 /usr/local/bin/sync.py
